# conf/system/systemd/system/vpn-watcher.service
[Unit]
Description=VPN Watcher
After=openvpn-client@bossk.service

[Service]
ExecStart=/home/wrk/check_vpn.sh

[Install]
WantedBy=multi-user.target
