# conf/system/systemd/system/dnsmasq.service.d/start.conf
[Unit]
After=network.target systemd-networkd.service
